# sol-bot/systemd/sol-threads-analytics.timer
[Unit]
Description=Run Sol Threads analytics sync hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
AccuracySec=5min
Persistent=true
Unit=sol-threads-analytics.service

[Install]
WantedBy=timers.target
