# sol-bot/systemd/sol-threads-analytics.service
[Unit]
Description=Sol Threads analytics sync
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/x-bot/sol-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /root/x-bot/sol-bot/threads_analytics.py sync --limit 50
